# systemd-units/chromebot-irc.socket
[Unit]
Description=chromebot IRC FIFO

[Socket]
ListenFIFO=/home/chromebot/chromebot-irc.fifo
